# web/deploy/systemd/material-pull-backend-dev.service
[Unit]
Description=Material Pull System Backend (Dev)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@RUN_USER@
Group=@RUN_USER@
WorkingDirectory=@PROJECT_ROOT@/backend
EnvironmentFile=-@PROJECT_ROOT@/deploy/scripts/material-pull.env
Environment=JAVA_HOME=@JAVA_HOME@
Environment=PATH=@JAVA_HOME@/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/usr/bin/mvn -f @PROJECT_ROOT@/backend/pom.xml clean package -DskipTests -q
ExecStart=@JAVA_HOME@/bin/java -jar @PROJECT_ROOT@/backend/target/@BACKEND_JAR@ --spring.profiles.active=@SPRING_PROFILE@
Restart=on-failure
RestartSec=10
StandardOutput=append:@LOG_DIR@/backend.log
StandardError=append:@LOG_DIR@/backend.log

[Install]
WantedBy=multi-user.target
